# roles/bacula-client/templates/bacula-fd2-encrypted.service.j2
.include /usr/lib/systemd/system/bacula-fd.service

[Service]
Environment=CONFIG={{ bacula_encrypted_config_file }}
